# python/ansible/hashicorp/collections/ansible_collections/enmanuelmoreira/hashicorp/roles/packer/templates/etc/profile.d/packer.sh.j2
# {{ ansible_managed }}
export PACKER_LOG="{{ packer_activate_log }}"
export PACKER_LOG_PATH="{{ packer_log_path }}"
